# daemons/go-worker.conf
description "Service for a nsqd"
author      "Kristina Kovalevskaya <isitiriss@gmail.com>"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [!2345]

# Automatically restart process if crashed
respawn

#pre-start script
#   mkdir -p /var/log/go-Nominatim/
#end script

# Essentially lets upstart know the process will detach itself to the background
expect fork

# Start the process
script
CURR_PATH="$(pwd)"
PATH_WORKER="/home/k/work/go/src/Nominatim/worker"
cd "${PATH_WORKER}"
exec ./worker
#cd "$(CURR_PATH)"
end script

#script
#su -$USER -c "exec /home/k/work/go/src/Nominatim/worker 2>&1" >>  $HOME/go-nominatim/worker.log
#end script
